# backup/dock2box-backup.timer
[Unit]
Description=Backup Dock2Box every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=dock2box-backup.service
