# backup/dock2box-backup.service
[Unit]
Description=Backup Dock2Box
Requires=dock2box.service

[Service]
Type=simple
ExecStart=/usr/local/bin/dock2box-backup.sh
