# boot-images/alpine3.4/dock2box
#!/sbin/openrc-run

description="Start install of operating system"

retry=300 # 5 min.

depend()
{
        use network dns hwdrivers sshd
}

start()
{
        ebegin "Starting Dock2Box install"
        cd /root
        ./install.sh
        eend $?
}

