# agent/files/ebuild/sys-apps/peekaboo/files/peekaboo.init
#!/sbin/runscript
# $Header: $

name=peekaboo
pidfile="/var/run/peekaboo.pid"
command="/usr/bin/peekaboo"
# required because daemon does not create pidfile or go into the background nicely
start_stop_daemon_args="-b -m -p ${pidfile}"
command_args="-s ${PEEKABOO_STATIC} -t ${PEEKABOO_TEMPLATES}"

depend() {
    use net
}

